# packaging/siphon.service
[Unit]
Description=siphon - High-performance DLP scanner REST API
Documentation=https://github.com/oxide11/siphon-rs
After=network.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/siphon serve
Restart=on-failure
RestartSec=5s

# Run as a dedicated unprivileged user. The .deb / .rpm post-install
# script creates the siphon user.
User=siphon
Group=siphon

# Pull configuration from /etc/siphon/siphon.env if present.
EnvironmentFile=-/etc/siphon/siphon.env

# Default state and audit log directories.
StateDirectory=siphon
LogsDirectory=siphon
ConfigurationDirectory=siphon

# Hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/siphon /var/log/siphon
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount @debug @cpu-emulation @obsolete
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0027

# Resource limits
LimitNOFILE=65536
TasksMax=4096

[Install]
WantedBy=multi-user.target
